# pi/systemd/home-energy-annual-analysis.timer
[Unit]
Description=Daily Home Energy annual analysis refresh

[Timer]
OnCalendar=03:30
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
