# pi/systemd/home-energy-supabase-forward-sync.timer
[Unit]
Description=Daily Home Energy Supabase forward sync

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true

[Install]
WantedBy=timers.target
